# quadlet/npm.container
[Unit]
Description=Nginx Proxy Manager Container
After=network-online.target

[Container]
ContainerName=npm
Image=docker.io/jc21/nginx-proxy-manager:latest

# Ports
PublishPort=80:80
PublishPort=81:81
PublishPort=443:443

# Volumes dengan label SELinux (:Z)
Volume=npm_letsencrypt:/etc/letsencrypt:Z
Volume=npm_data:/data:Z

# --- PENYEMPURNAAN HEALTHCHECK ---
# Gunakan port 81 (dashboard admin) untuk cek karena itu yang paling cepat merespon
HealthCmd=curl -f -s http://localhost:81 || exit 1
HealthInterval=30s
HealthRetries=3
# NPM butuh waktu untuk startup database internalnya
HealthStartPeriod=20s

LogDriver=journald

[Service]
# Otomatis restart jika crash
Restart=always
# Berikan waktu jeda sebelum restart agar tidak loop cepat jika DB bermasalah
RestartSec=5s

[Install]
WantedBy=multi-user.target
